# app/supervisord.conf
[program:nginx]
command=/usr/sbin/nginx

[program:reids]
command=/usr/bin/reids-server

[program:myuwsgi]
command=/usr/local/bin/uwsgi --ini /root/app/uwsgi.ini
